# backup-kcal.timer
[Unit]
Description=Repeats backup-kcal.service daily.

[Timer]
OnCalendar=daily
Unit=backup-kcal.service

[Install]
WantedBy=multi-user.target
